# power-management/addons/turbopm/turbopm.service
[Unit]
Description=launcher for sys76pm service             
[Service]
Type=simple  
ExecStart=/usr/local/bin/turbopm
[Install]
WantedBy=default.target
