# power-management/addons/ppd/ppdpm.service
[Unit]
Description=launcher for ppdpm service             
[Service]
Type=simple  
ExecStart=/usr/local/bin/ppdpm
[Install]
WantedBy=default.target
